# ConfidentialComputing/CvmImageBuilder/rootfs-files/etc-overlay.service
[Unit]
Description=Setup writable /etc overlay on read-only root
DefaultDependencies=no
Conflicts=shutdown.target
Before=local-fs.target cloud-init-local.service cloud-init.service
After=systemd-remount-fs.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/setup-etc-overlay.sh

[Install]
WantedBy=local-fs.target
